# supervisord.conf
[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[supervisord]
logfile=/var/log/supervisor/supervisor.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=true                                   ; (start in foreground if true;default false)
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/supervisor/               ; where child log files will live
;umask=022                  ; (process file creation umask;default 022)
;identifier=supervisor       ; (supervisord identifier, default is 'supervisor')
;directory=/tmp              ; (default is not to cd during start)
;nocleanup=true              ; (don't clean up tempfiles at start;default false)
;childlogdir=/tmp            ; ('AUTO' child log dir, default $TEMP)
;environment=KEY=value       ; (key value pairs to add to environment)
;strip_ansi=false            ; (strip ansi escape codes in logs; default false)
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)

[program:rsyslogd]
command=/usr/sbin/rsyslogd -n
process_name=%(program_name)s
autostart=true
autorestart=false
user=root
directory=/
priority=912
stdout_logfile=/var/log/%(program_name)s-stdout.log
stderr_logfile=/var/log/%(program_name)s-stderr.log

[program:cron]
command = /usr/sbin/cron -f
stdout_logfile = /var/log/%(program_name)s-stdout.log
stderr_logfile = /var/log/%(program_name)s-stderr.log
autostart=true
autorestart = false

[program:bind9]
autostart=true
autorestart = false
command=/etc/init.d/bind9 restart

[program:apache2]
autostart=true
autorestart = false
command=/etc/init.d/apache2 restart

[program:fail2ban]
autostart=true
autorestart = false
command = /etc/init.d/fail2ban restart

[program:denyhosts]
autostart=true
autorestart = false
command = /etc/init.d/denyhosts restart

[program:mailman]
autostart=true
autorestart = false
command=/usr/lib/mailman/bin/mailmanctl restart

[program:courier-imap]
command=/etc/init.d/courier-imap restart
priority=3
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart = true

[program:courier-imap-ssl]
command=/etc/init.d/courier-imap-ssl restart
priority=3
autostart = true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:courier-authdaemon]
command=/etc/init.d/courier-authdaemon restart
priority=2
autostart = true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:saslauthd]
command=/etc/init.d/saslauthd restart
priority=1
autostart = true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:clamav]
command=/etc/init.d/clamav-daemon restart
priority=2
autostart = true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:freshclam]
command=/etc/init.d/clamav-freshclam restart
priority=2
autostart = true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:spamassasin]
command=/etc/init.d/spamassassin restart
priority=2
autostart = true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:amavis]
process_name = amavis
command = /etc/init.d/amavis restart
startsecs = 0
autorestart = false
priority = 20

[program:postfix]
process_name = postfix
command = /etc/init.d/postfix restart
startsecs = 0
autorestart = false
priority = 20

[program:postgrey]
process_name = postgrey
command = /etc/init.d/postgrey start
startsecs = 0
autorestart = false
priority = 20


[program:haveged]
command=/etc/init.d/haveged restart
priority=3
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart = true






